# workspace/ttnhabbridge/src/dist/ttnhabbridge.service
# copy this to /etc/systemd/system
# then sudo systemctl daemon-reload
# and enable with sudo systemctl enable --now ttnhabbridge

[Unit]
Description=ttnhabbridge
After=network.target

[Service]
WorkingDirectory=/home/bertrik/ttnhabbridge
ExecStart=/home/bertrik/ttnhabbridge/bin/ttnhabbridge
Type=simple
KillMode=process
Restart=on-failure
RestartSec=10s
TimeoutStopSec=15s
User=bertrik

[Install]
WantedBy=multi-user.target
